# vps-receiver.service
[Unit]
Description=VPS Results Receiver API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/vps-receiver
Environment="PATH=/home/ubuntu/vps-receiver/venv/bin"
Environment="VPS_SECRET_KEY=votre_cle_secrete_ici"
ExecStart=/home/ubuntu/vps-receiver/venv/bin/gunicorn -w 4 -b 0.0.0.0:5000 vps_receiver:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target